# .config/cron/crontab
*/15 * * * * cron-rss
*/15 * * * * cron-pacupdates
*/3 * * * * mailsync
